# home/.config/autostart/vpn-serve.desktop
[Desktop Entry]
Exec=/usr/bin/vpn-serve
Icon=dialog-scripts
Name=vpn-serve
Path=
Type=Application
X-KDE-AutostartScript=true
